# src/usr/local/etc/rc.d/cloudflared
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: cloudflared
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.d/cloudflared to enable cloudflared:
#
# cloudflared_enable (bool):	Set to "NO" by default.
#				Set it to "YES" to enable cloudflared.
# cloudflared_token_file (str):	Path to the token file.
# cloudflared_quic_disable_pmtu_discovery (bool): Set to "YES" by default.

. /etc/rc.subr

name="cloudflared"
rcvar="cloudflared_enable"

load_rc_config $name

: ${cloudflared_enable:="NO"}
: ${cloudflared_token_file:="/usr/local/etc/cloudflared/token"}
: ${cloudflared_quic_disable_pmtu_discovery:="YES"}

command="/usr/local/bin/cloudflared"
pidfile="/var/run/${name}.pid"
procname="/usr/local/bin/cloudflared"

# Build arguments based on configuration
command_args="tunnel run --token-file ${cloudflared_token_file}"

if checkyesno cloudflared_quic_disable_pmtu_discovery; then
    command_args="${command_args} --quic-disable-pmtu-discovery"
fi

# We use daemon to manage the background process and PID
actual_command="${command}"
command="/usr/sbin/daemon"
command_args="-c -f -p ${pidfile} ${actual_command} ${command_args}"

run_rc_command "$1"
